# ASO/cargo-ios/__fat-deb-root/Library/LaunchDaemons/com.nkaso.cargo.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nkaso.cargo.daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/cargodaemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ExecuteAllowed</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/var/mobile/Media/cargodaemon.log</string>
	<key>StandardOutPath</key>
	<string>/private/var/mobile/Media/cargodaemon.log</string>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
